# deployments/webplow.service
[Unit]
Description=Webplow - WebP conversion API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/webplow
ExecStart=/opt/webplow/webplow
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
EnvironmentFile=/opt/webplow/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/imgproxy/uploads /opt/webplow
PrivateTmp=true

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
